# keystone_idm.conf
# keystone_idm - keystone_idm job file
description "Keystone server extended to use OAuth2.0"
author "Enrique G. Navalon <garcianavalon@gmail.com>"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
# Automatically restart process if crashed
respawn
setuid root
script
cd /home/isabel/idm/keystone/
#activate the venv
. .venv/bin/activate
#run keystone
bin/keystone-all
end script
            