# doc/example_configs/keystone_idm.conf
# keystone_idm - keystone_idm job file
description "Service conf file for the IdM backend based in Keystone"
author "Enrique Garcia Navalon <garcianavalon@gmail.com>"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
# Automatically restart process if crashed
respawn
setuid root
script
cd $absolute_keystone_path
#activate the venv
. .venv/bin/activate
#run keystone
bin/keystone-all
end script
            